# supervisord.conf
[program:dataprep]
command=java -Xmx8g -jar /spv2-dataprep/dataprep.jar
priority=10

[program:db_worker]
command=python3 server.py
priority=20
