# supervisord.conf
[supervisord]

[program:spv2-dataprep]
command=java -Xmx8g -jar /spv2-dataprep/dataprep.jar
priority=10

[program:spv2-server]
command=python3 server.py
priority=20
